# docs/com.roughdiamond.chrome-debug.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.roughdiamond.chrome-debug</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Google Chrome.app/Contents/MacOS/Google Chrome</string>
    <string>--remote-debugging-port=9223</string>
    <string>--no-first-run</string>
    <string>--no-default-browser-check</string>
    <string>--user-data-dir=/Users/judy/Library/Application Support/Google/Chrome/RoughDiamond</string>
    <string>https://www.coupang.com</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/chrome-debug.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/chrome-debug-error.log</string>
</dict>
</plist>
